# files/etc/uci-defaults/01_force-lan-dhcp
#!/bin/sh

# set force lan to dhcp.

uci set dhcp.lan.force='1'
uci commit dhcp

exit 0