# files/etc/uci-defaults/02_disable_peerdns
#!/bin/sh

# this fix dns leaks.
uci set network.wan.peerdns=0
uci commit network

exit 0