# files/etc/uci-defaults/02_set-classless-routes
#!/bin/sh

# this fixes the tunnelvision/tunnelcrack vulnerability for wireguard

uci set network.wan.classlessroute='0'

uci commit network

exit 0